# features/github_action_runner/file.include/usr/lib/dracut/modules.d/98tar-rootfs/sysroot.mount
[Unit]
Before=initrd-root-fs.target
DefaultDependencies=no

[Mount]
What=tmpfs
Where=/sysroot
Type=tmpfs
