# features/orabos/file.include/etc/systemd/system/ovsdb-server.service.d/allow-group.conf
[Service]
UMask=0002
RuntimeDirectoryMode=0775

