# features/orabos/file.include/etc/systemd/system/ovsdb-server.service.d/ovsdb-server-permissions.conf
[Service]
ExecStartPost=-/usr/bin/chmod g+rw /run/openvswitch/db.sock
ExecStartPost=-/usr/bin/chgrp openvswitch /run/openvswitch/db.sock
