# features/orabos/file.include/etc/systemd/system/ovsdb-server.service.d/ovsdb-server-chmod.conf
[Service]
ExecStartPost=-/usr/bin/chmod g+rw /run/openvswitch/db.sock
